# deploy/proxmox/forumline-fleet-sync.service
[Unit]
Description=Forumline fleet sync — ensure all LXCs have deploy key + syslog config

[Service]
Type=oneshot
ExecStart=/etc/forumline/forumline-fleet-sync.sh
